# systemd/autokernel-continuous.service
[Unit]
Description=AutoKernel Continuous Optimization Pipeline
After=network.target ollama.service
Wants=ollama.service

[Service]
Type=simple
User=alexendros
WorkingDirectory=/home/alexendros/repositorios/org-iniciativas-alexendros/autokernel
Environment=PATH=/home/alexendros/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=HOME=/home/alexendros
ExecStartPre=/bin/sh -c 'command -v uv >/dev/null || { echo "uv no encontrado"; exit 1; }'
ExecStart=/home/alexendros/repositorios/org-iniciativas-alexendros/autokernel/scripts/continuous_pipeline.py
StandardOutput=append:/home/alexendros/logs/autokernel-continuous.log
StandardError=inherit
Restart=on-failure
RestartSec=30
TimeoutSec=0
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7

[Install]
WantedBy=multi-user.target
