# systemd/autokernel-nightly.timer
[Unit]
Description=Run AutoKernel nightly optimization at 2AM

[Timer]
OnCalendar=*-*-* 02:00:00
Persistent=true
RandomizedDelaySec=900

[Install]
WantedBy=timers.target
